# lab_doc/data/inittab
::sysinit:/etc/init.d/rcS

::respawn:/sbin/getty -L ttyAMA0 115200 xterm

